# auto_login.timer
[Unit]
Description = Run Auto_Login Every Hour

[Timer]
OnUnitActiveSec = 1h
Unit = auto_login.service

[Install]
WantedBy = multi-user.target
